# gasmon/gasmon.cron
# $Id: gasmon.cron,v 1.1 2007-02-15 15:52:47 jschamba Exp $
# $Log: not supported by cvs2svn $
# Revision 1.2  2003/09/30 20:41:50  geurts
# Added CVS keywords
#
#
# -- TOFr GasMonitor Cron Job
#
* * * * * $HOME/bin/gasmon -readout >> /data/tofr/gasmon.log
